# data/enriched_alerts.json
[
  {"alert_id": "ALT-10001", "rule_description": "Suspicious PowerShell Encoded Command Execution", "level": 15, "technique_id": "T1059.001", "agent_name": "WIN-DC01"},
  {"alert_id": "ALT-10002", "rule_description": "Mimikatz Credential Dumping Detected", "level": 15, "technique_id": "T1003", "agent_name": "WIN-DC01"},
  {"alert_id": "ALT-10003", "rule_description": "New Scheduled Task Created for Persistence", "level": 10, "technique_id": "T1053.005", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10004", "rule_description": "Outbound Connection to Known C2 Infrastructure", "level": 15, "technique_id": "T1071.001", "agent_name": "LNX-WEB01"},
  {"alert_id": "ALT-10005", "rule_description": "Registry Run Key Modified for Autostart", "level": 8, "technique_id": "T1547.001", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10006", "rule_description": "Suspicious Use of Rundll32 to Load DLL", "level": 12, "technique_id": "T1218.011", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10007", "rule_description": "Brute Force Attempt on SSH Service", "level": 9, "technique_id": "T1110", "agent_name": "LNX-WEB01"},
  {"alert_id": "ALT-10008", "rule_description": "PsExec Lateral Movement Detected", "level": 13, "technique_id": "T1021.002", "agent_name": "WIN-DC01"},
  {"alert_id": "ALT-10009", "rule_description": "Disabling of Windows Defender Real-Time Protection", "level": 14, "technique_id": "T1562.001", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10010", "rule_description": "Unusual Volume of DNS Requests (Possible Tunneling)", "level": 11, "technique_id": "T1071.004", "agent_name": "LNX-DB01"},
  {"alert_id": "ALT-10011", "rule_description": "Suspicious Process Injection into Explorer.exe", "level": 13, "technique_id": "T1055", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10012", "rule_description": "New Local Administrator Account Created", "level": 10, "technique_id": "T1136.001", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10013", "rule_description": "Archive File Created in Temp Directory (Possible Staging)", "level": 6, "technique_id": "T1560.001", "agent_name": "LNX-DB01"},
  {"alert_id": "ALT-10014", "rule_description": "Failed Login Followed by Successful Login (Anomaly)", "level": 7, "technique_id": "T1110.001", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10015", "rule_description": "Suspicious Use of Certutil to Download File", "level": 12, "technique_id": "T1105", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10016", "rule_description": "Kerberoasting Activity Detected", "level": 15, "technique_id": "T1558.003", "agent_name": "WIN-DC01"},
  {"alert_id": "ALT-10017", "rule_description": "Web Shell Upload Detected on Web Server", "level": 14, "technique_id": "T1505.003", "agent_name": "LNX-WEB01"},
  {"alert_id": "ALT-10018", "rule_description": "Clearing of Windows Event Logs", "level": 13, "technique_id": "T1070.001", "agent_name": "WIN-DC01"},
  {"alert_id": "ALT-10019", "rule_description": "Unusual Outbound Data Transfer Volume", "level": 9, "technique_id": "T1041", "agent_name": "LNX-DB01"},
  {"alert_id": "ALT-10020", "rule_description": "Suspicious WMI Event Subscription Persistence", "level": 11, "technique_id": "T1546.003", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10021", "rule_description": "Multiple Failed Sudo Attempts Detected", "level": 5, "technique_id": "T1548.003", "agent_name": "LNX-WEB01"},
  {"alert_id": "ALT-10022", "rule_description": "Suspicious Named Pipe Impersonation", "level": 12, "technique_id": "T1134", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10023", "rule_description": "Port Scan Detected from Internal Host", "level": 6, "technique_id": "T1046", "agent_name": "LNX-DB01"},
  {"alert_id": "ALT-10024", "rule_description": "Suspicious Office Macro Execution Spawning Shell", "level": 14, "technique_id": "T1204.002", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10025", "rule_description": "Unsigned Binary Executed from Temp Path", "level": 8, "technique_id": "T1027", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10026", "rule_description": "Golden Ticket Indicator - Anomalous TGT Lifetime", "level": 15, "technique_id": "T1558.001", "agent_name": "WIN-DC01"},
  {"alert_id": "ALT-10027", "rule_description": "Suspicious cron Job Modification", "level": 9, "technique_id": "T1053.003", "agent_name": "LNX-WEB01"},
  {"alert_id": "ALT-10028", "rule_description": "Exfiltration Over Alternative Protocol (FTP)", "level": 12, "technique_id": "T1048", "agent_name": "LNX-DB01"},
  {"alert_id": "ALT-10029", "rule_description": "Suspicious LSASS Memory Access", "level": 15, "technique_id": "T1003.001", "agent_name": "WIN-DC01"},
  {"alert_id": "ALT-10030", "rule_description": "Antivirus Signature Update Failure", "level": 3, "technique_id": "T1562.001", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10031", "rule_description": "New Firewall Rule Added Allowing Inbound RDP", "level": 10, "technique_id": "T1562.004", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10032", "rule_description": "Suspicious PowerShell Download Cradle", "level": 13, "technique_id": "T1059.001", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10033", "rule_description": "Reconnaissance Command Execution (whoami/net)", "level": 4, "technique_id": "T1082", "agent_name": "LNX-WEB01"},
  {"alert_id": "ALT-10034", "rule_description": "Unusual Time-of-Day Login Activity", "level": 5, "technique_id": "T1078", "agent_name": "LNX-DB01"},
  {"alert_id": "ALT-10035", "rule_description": "Suspicious DLL Sideloading Detected", "level": 12, "technique_id": "T1574.002", "agent_name": "WIN-APP02"},
  {"alert_id": "ALT-10036", "rule_description": "Ransomware Note File Created", "level": 15, "technique_id": "T1486", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10037", "rule_description": "Suspicious VSS Admin Shadow Copy Deletion", "level": 15, "technique_id": "T1490", "agent_name": "WIN-FIN03"},
  {"alert_id": "ALT-10038", "rule_description": "SSH Login from Unusual Geolocation", "level": 8, "technique_id": "T1078.003", "agent_name": "LNX-WEB01"},
  {"alert_id": "ALT-10039", "rule_description": "Suspicious Base64 Encoded Command in Bash History", "level": 11, "technique_id": "T1140", "agent_name": "LNX-DB01"},
  {"alert_id": "ALT-10040", "rule_description": "Unusual Volume of Failed Authentication Events", "level": 7, "technique_id": "T1110", "agent_name": "WIN-DC01"}
]
